# udev/050_tcos_devices.rules
# This rule exec a script to save some udev env vars in a tmp file
#

SUBSYSTEM=="block",                    RUN+="/sbin/tcos-udevd.sh"

