# debian/initramfs-tools-tcos.postrm
#!/bin/bash

set -e

if [ "x${1}" = "xpurge" ]; then
  update-rc.d tcos-standalone remove >/dev/null
fi



#DEBHELPER#

