# debian/tcos-standalone.postrm
#!/bin/bash

set -e

# on purge delete tcos-standalone daemon
case "$1" in
    purge|remove)
       update-rc.d -f tcos-standalone remove >/dev/null ||true
       rm -rf /var/lib/tcos/standalone/*
    ;;
esac



#DEBHELPER#

