# debian/tcos-standalone.prerm
#!/bin/bash
set -e

case "$1" in
    purge|remove)
  	# stop tcos-standalone script
  	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
     		invoke-rc.d tcos-standalone stop || true
  	else
   		 /etc/init.d/tcos-standalone stop || true
  	fi
    ;;
esac



#DEBHELPER#
